# app/movie.conf
[program:movie]
command=/root/myenv/bin/python /root/movie.py
directory=/root
autostart=true
autorestart=true
stderr_logfile=/var/log/movie.err.log
stdout_logfile=/var/log/movie.out.log
